# deploy/systemd/reblas-dashboard-live-watchdog.user.service
[Unit]
Description=Reblas Dashboard Live Web Watchdog (user)

[Service]
Type=oneshot
ExecStart=/projects/reblas-crew-dashboard/scripts/check-live-web.sh

